# deploy/systemd/myasset-stack.target
[Unit]
Description=MyAsset Application Stack (API + Nginx)
Requires=myasset-api.service nginx.service
After=myasset-api.service nginx.service

[Install]
WantedBy=multi-user.target

